# init/ethereum-delegated-tx-mainnet.service
[Unit]
Description=Ethereum Delegated Transactions Service for Mainnet

[Service]
Type=simple
Environment="MONGODB_URL=mongodb://localhost:27017"
Restart=always
RestartSec=1s
ExecStart=/usr/bin/npm run start
WorkingDirectory=/home/zitro/apps/ethereum-delegated-tx-service-mainnet

[Install]
WantedBy=multi-user.target
